# deployment/mcpx.service
[Unit]
Description=MCPX HTTP MCP Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/data/web/mcpx.lol

# Load environment variables from file
EnvironmentFile=/data/web/mcpx.lol/.env

# Set uv cache directory to avoid permission issues with /var/www/.cache
Environment="UV_CACHE_DIR=/data/web/mcpx.lol/data/.uv-cache"

# Run the server using uv
ExecStart=/usr/bin/uv run python -m mcp_server.main

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/web/mcpx.lol/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcpx

[Install]
WantedBy=multi-user.target

